# openssh-akcenv/openssh-akcenv-6.6p1/debian/openssh-akcenv.postrm
#!/bin/sh
set -eu

case $1 in
  purge)
    update-rc.d ssh remove >/dev/null || true

    if which deluser >/dev/null 2>&1; then
      deluser --quiet sshd > /dev/null || true
    fi
    ;;
esac
