# openssh-akcenv/openssh-akcenv-6.6p1/debian/openssh-akcenv.postrm
#!/bin/sh
set -eu
update-rc.d ssh remove >/dev/null || true
