# openssh-akcenv/openssh-akcenv-6.6p1/debian/openssh-akcenv.postinst
#!/bin/sh
set -eu
update-rc.d ssh defaults >/dev/null || true
